# fauxmo.conf
#!upstart
description "Fauxmo"

start on filesystem
stop on runlevel [06]

oom score 1000
exec /usr/local/bin/fauxmo.py
respawn
respawn limit 10 90
